# scripts/debian/bitcore-bitcoind.conf
description "Bitcoin Core for Bcccore"
author "Open Wallet Stack"

limit nofile 20000 30000

start on runlevel [2345]
stop on runlevel [016]

kill timeout 300
kill signal SIGINT

# user/group for bcccore daemon to run as
setuid bcccore
setgid bcccore

# home dir of the bcccore daemon user
env HOME=/home/bcccore

respawn
respawn limit 5 15

script
  exec bcccore-bcccored -datadir=/home/bcccore/.bcccore/data/
end script
